# Fase 2 - Analisis semantico NLP/mitre_mapping.json
{
  "CobaltStrike": ["T1059.001", "T1043", "T1087.001"],
  "Cobalt Strike": ["T1059.001", "T1043", "T1087.001"],
  "cobaltstrike": ["T1059.001", "T1043", "T1087.001"],
  "Mimikatz": ["T1003.001", "T1003.002"],
  "Metasploit": ["T1190", "T1210"],
  "Emotet": ["T1566.001", "T1059.005"],
  "Ryuk": ["T1486", "T1027"],
  "phishing": ["T1566"],
  "ransomware": ["T1486"],
  "lateral movement": ["T1021.001", "T1021.002"],
  "credential dumping": ["T1003"],
  "exploit": ["T1190"],
  "malware": ["T1204.002"],
  "backdoor": ["T1505.003"],
  "keylogger": ["T1056.001"],
  "spyware": ["T1552.001"],
  "trojan": ["T1204.002"],
  "worm": ["T1204.002"],
  "botnet": ["T1583.001"],
  "DDoS": ["T1498"],
  "SQL injection": ["T1190"],
  "XSS": ["T1190"],
  "zero-day": ["T1190"],
  "APT": ["T1027", "T1082"],
  "C2": ["T1071.001", "T1573.001"],
  "payload": ["T1204.002"],
  "exfiltration": ["T1041"],
  "encryption": ["T1486"],
  "decryption": ["T1486"],
  "privilege escalation": ["T1068"],
  "persistence": ["T1037", "T1547.001"],
  "evade detection": ["T1027", "T1089"],
  "command and control": ["T1071", "T1573"],
  "data theft": ["T1020", "T1041"],
  "network scanning": ["T1046"],
  "vulnerability scanning": ["T1046"],
  "brute force": ["T1110"],
  "social engineering": ["T1566"],
  "man-in-the-middle": ["T1557"],
  "rootkit": ["T1014"],
  "bootkit": ["T1542.003"],
  "cryptojacking": ["T1496"],
  "fileless": ["T1059.001", "T1059.005"],
  "polymorphic": ["T1027"],
  "obfuscation": ["T1027"],
  "steganography": ["T1027"],
  "registry manipulation": ["T1112"],
  "process injection": ["T1055"],
  "hook": ["T1055"],
  "keystroke logging": ["T1056.001"],
  "screen capture": ["T1113"],
  "audio capture": ["T1123"],
  "video capture": ["T1125"],
  "input capture": ["T1056"],
  "credential access": ["T1003", "T1555"],
  "token manipulation": ["T1134"],
  "session hijacking": ["T1563.001"],
  "cookie theft": ["T1539"],
  "browser exploitation": ["T1189"],
  "drive-by compromise": ["T1189"],
  "watering hole": ["T1189"],
  "supply chain attack": ["T1195"],
  "living off the land": ["T1218"],
  "dual use tools": ["T1218"],
  "scripting": ["T1059"],
  "PowerShell": ["T1059.001"],
  "Python": ["T1059.006"],
  "Bash": ["T1059.004"],
  "JavaScript": ["T1059.007"],
  "VBScript": ["T1059.005"],
  "macro": ["T1059.005"],
  "executable": ["T1204.002"],
  "DLL": ["T1204.002"],
  "service": ["T1543.003"],
  "scheduled task": ["T1053.005"],
  "startup": ["T1547.001"],
  "registry run": ["T1547.001"],
  "WMI": ["T1047"],
  "RDP": ["T1021.001"],
  "SSH": ["T1021.004"],
  "VPN": ["T1071.001"],
  "Tor": ["T1071.001"],
  "proxy": ["T1090"],
  "tunnel": ["T1572"],
  "DNS tunneling": ["T1071.004"],
  "ICMP tunneling": ["T1071.003"],
  "HTTP tunneling": ["T1071.001"],
  "HTTPS tunneling": ["T1071.001"],
  "SMB": ["T1021.002"],
  "FTP": ["T1071.002"],
  "SMTP": ["T1071.003"],
  "IMAP": ["T1071.003"],
  "POP3": ["T1071.003"],
  "RPC": ["T1021.002"],
  "LDAP": ["T1087.002"],
  "Kerberos": ["T1558.003"],
  "NTLM": ["T1558.002"],
  "Active Directory": ["T1087.002"],
  "Group Policy": ["T1484.001"],

  "LockBit": ["T1486", "T1027"],
  "LockBit 3.0": ["T1486", "T1027"],
  "Conti": ["T1486", "T1027"],
  "REvil": ["T1486", "T1027"],
  "BlackCat": ["T1486", "T1027"],
  "DarkSide": ["T1486", "T1027"],
  "Hive": ["T1486", "T1027"],
  "BlackBasta": ["T1486", "T1027"],
  "Royal": ["T1488", "T1027"],
  "Play": ["T1486", "T1027"],
  "Cl0p": ["T1486", "T1195"],
  "Maze": ["T1486", "T1041"],
  "Dharma": ["T1486", "T1059"],
  "GandCrab": ["T1486", "T1059"],
  "Phobos": ["T1486", "T1059"],
  "Babuk": ["T1486", "T1027"],
  "Avaddon": ["T1486", "T1566"],
  "NetWalker": ["T1486", "T1059"],
  "Revil": ["T1486", "T1027"],

  "Sliver": ["T1059.001", "T1071.001"],
  "Brute Ratel": ["T1059.001", "T1071.001"],
  "BruteRatel": ["T1059.001", "T1071.001"],
  "Havoc": ["T1059.001", "T1071.001"],
  "Mythic": ["T1059.001", "T1071.001"],
  "Silver": ["T1059.001", "T1071.001"],
  "Empire": ["T1059.001", "T1071.001"],
  "Covenant": ["T1059.001", "T1071.001"],
  "PoshC2": ["T1059.001", "T1071.001"],
  "Nighthawk": ["T1059.001", "T1071.001"],
  "Rubeus": ["T1558.003", "T1558.002"],
  "BloodHound": ["T1087.002", "T1018"],
  "SharpHound": ["T1087.002", "T1018"],
  "Responder": ["T1557", "T1040"],
  "Impacket": ["T1021.002", "T1557"],
  "CrackMapExec": ["T1021.002", "T1087.002"],
  "PsExec": ["T1021.002", "T1569.002"],
  "WinPEAS": ["T1068", "T1082"],
  "LinPEAS": ["T1068", "T1082"],
  "LaZagne": ["T1003", "T1555"],
  "Seatbelt": ["T1082", "T1033"],
  "SharpUp": ["T1068", "T1082"],

  "Nmap": ["T1046"],
  "Masscan": ["T1046"],
  "Zmap": ["T1046"],
  "Shodan": ["T1046", "T1596"],
  "Censys": ["T1046", "T1596"],
  "theHarvester": ["T1593", "T1594"],
  "Recon-ng": ["T1593", "T1594"],
  "Maltego": ["T1593", "T1594"],
  "SpiderFoot": ["T1593", "T1594"],
  "Sublist3r": ["T1596.005"],
  "Amass": ["T1596.005"],
  "Gobuster": ["T1596.005"],
  "Dirb": ["T1596.005"],
  "Dirbuster": ["T1596.005"],

  "Log4Shell": ["T1190"],
  "Log4j": ["T1190"],
  "EternalBlue": ["T1210"],
  "ProxyLogon": ["T1190"],
  "ProxyShell": ["T1190"],
  "ProxyNotShell": ["T1190"],
  "PrintNightmare": ["T1190", "T1068"],
  "BlueKeep": ["T1210"],
  "Zerologon": ["T1190"],
  "Dirty Pipe": ["T1068"],
  "DirtyCow": ["T1068"],
  "Heartbleed": ["T1040"],
  "ShellShock": ["T1190"],
  "Stagefright": ["T1203"],

  "carding": ["T1020", "T1041"],
  "fullz": ["T1005", "T1039"],
  "dumps": ["T1005", "T1039"],
  "skimmer": ["T1056.001", "T1041"],
  "ATM": ["T1486", "T1005"],
  "SWIFT": ["T1486", "T1020"],
  "wire fraud": ["T1486"],
  "money laundering": ["T1486", "T1020"],
  "darknet market": ["T1583.001"],
  "deep web": ["T1583.001"],
  "onion site": ["T1071.001"],
  "forum": ["T1583.001"],
  "marketplace": ["T1583.001"],

  "RAT": ["T1059.001", "T1071.001"],
  "remote access": ["T1021.001", "T1071.001"],
  "webshell": ["T1505.003", "T1569.002"],
  "beacon": ["T1071.001"],
  "implant": ["T1505.003", "T1059.001"],
  "dropper": ["T1204.002", "T1059.001"],
  "loader": ["T1204.002", "T1059.001"],
  "stealer": ["T1005", "T1555"],
  "infostealer": ["T1005", "T1555"],
  "formgrabber": ["T1056.001", "T1040"],
  "screenshot": ["T1113"],
  "clipboard": ["T1115"],
  "keystroke": ["T1056.001"], 
  "escalate privileges": ["T1068"],
  "elevation": ["T1068"],
  "sudo": ["T1548.001"],
  "SUID": ["T1548.001"],
  "UAC bypass": ["T1548.002"],
  "token impersonation": ["T1134"],
  "impersonate": ["T1134"],

  "data exfiltration": ["T1041"],
  "data staging": ["T1074"],
  "data compression": ["T1560"],
  "data encoding": ["T1132"],
  "data encryption": ["T1560.001"],
  "archive": ["T1560"],
  "base64": ["T1132.001"],
  "hex encoding": ["T1132.002"],
  "dns over https": ["T1071.004"],
  "DoH": ["T1071.004"],
  "domain fronting": ["T1090.004"],

  "persistence mechanism": ["T1037", "T1547"],
  "registry persistence": ["T1547.001", "T1112"],
  "scheduled persistence": ["T1053.005"],
  "boot persistence": ["T1547.001", "T1542"],
  "service persistence": ["T1543.003"],
  "WMI persistence": ["T1546.003"],

  "defense evasion": ["T1027", "T1562"],
  "anti-forensics": ["T1027", "T1070"],
  "timestomp": ["T1070.006"],
  "log clearing": ["T1070.001"],
  "indicator removal": ["T1070"],
  "file deletion": ["T1070.004"],
  "process hollowing": ["T1055.012"],
  "DLL injection": ["T1055.001"],
  "process doppelganging": ["T1055.013"],
  "AMSI bypass": ["T1562.001"],
  "ETW bypass": ["T1562.006"],

  "collection": ["T1005", "T1039"],
  "sensitive data": ["T1005", "T1039"],
  "credentials": ["T1003", "T1555"],
  "passwords": ["T1003", "T1555"],
  "hashes": ["T1003.001"],
  "NTDS": ["T1003.003"],
  "SAM database": ["T1003.002"],
  "LSASS": ["T1003.001"],
  "browser passwords": ["T1555.003"],
  "wifi passwords": ["T1557.001"],

  "lateral": ["T1021"],
  "pivot": ["T1021"],
  "jump box": ["T1021"],
  "psexec lateral": ["T1021.002"],
  "wmi lateral": ["T1047"],
  "dcom": ["T1021.003"],
  "smb lateral": ["T1021.002"],

  "credential harvesting": ["T1003", "T1110"],
  "password spray": ["T1110.003"],
  "password stuffing": ["T1110.004"],
  "credential stuffing": ["T1110.004"],
  "pass-the-hash": ["T1550.002"],
  "pass-the-ticket": ["T1550.003"],
  "golden ticket": ["T1558.001"],
  "silver ticket": ["T1558.002"],
  "kerberoasting": ["T1558.003"],
  "AS-REP roasting": ["T1558.004"],

  "data breach": ["T1041", "T1020"],
  "data leak": ["T1041", "T1020"],
  "data dump": ["T1005", "T1041"],
  "data leak site": ["T1041"],
  "leak": ["T1041", "T1005"],
  "steal data": ["T1005", "T1041"],
  "exfiltrate": ["T1041"],
  "siphon": ["T1041"],

  "supply chain": ["T1195"],
  "dependency confusion": ["T1195.002"],
  "compromised update": ["T1195.002"],
  "trojanized": ["T1195.002"],

  "container escape": ["T1611"],
  "kubernetes": ["T1610"],
  "docker": ["T1610"],
  "cloud": ["T1078"],
  "AWS": ["T1078.004"],
  "Azure": ["T1078.004"],
  "GCP": ["T1078.004"]
}