# Fase 2 - Analisis semantico NLP/mitre_mapping.json
{
  "CobaltStrike": ["T1059.001", "T1043", "T1087.001"],
  "Mimikatz": ["T1003.001", "T1003.002"],
  "Metasploit": ["T1190", "T1210"],
  "Emotet": ["T1566.001", "T1059.005"],
  "Ryuk": ["T1486", "T1027"],
  "phishing": ["T1566"],
  "ransomware": ["T1486"],
  "lateral movement": ["T1021.001", "T1021.002"],
  "credential dumping": ["T1003"],
  "exploit": ["T1190"],
  "malware": ["T1204.002"],
  "backdoor": ["T1505.003"],
  "keylogger": ["T1056.001"],
  "spyware": ["T1552.001"],
  "trojan": ["T1204.002"],
  "worm": ["T1204.002"],
  "botnet": ["T1583.001"],
  "DDoS": ["T1498"],
  "SQL injection": ["T1190"],
  "XSS": ["T1190"],
  "zero-day": ["T1190"],
  "APT": ["T1027", "T1082"],
  "C2": ["T1071.001", "T1573.001"],
  "payload": ["T1204.002"],
  "exfiltration": ["T1041"],
  "encryption": ["T1486"],
  "decryption": ["T1486"],
  "privilege escalation": ["T1068"],
  "persistence": ["T1037", "T1547.001"],
  "evade detection": ["T1027", "T1089"],
  "command and control": ["T1071", "T1573"],
  "data theft": ["T1020", "T1041"],
  "network scanning": ["T1046"],
  "vulnerability scanning": ["T1046"],
  "brute force": ["T1110"],
  "social engineering": ["T1566"],
  "man-in-the-middle": ["T1557"],
  "rootkit": ["T1014"],
  "bootkit": ["T1542.003"],
  "cryptojacking": ["T1496"],
  "fileless": ["T1059.001", "T1059.005"],
  "polymorphic": ["T1027"],
  "obfuscation": ["T1027"],
  "steganography": ["T1027"],
  "registry manipulation": ["T1112"],
  "process injection": ["T1055"],
  "hook": ["T1055"],
  "keystroke logging": ["T1056.001"],
  "screen capture": ["T1113"],
  "audio capture": ["T1123"],
  "video capture": ["T1125"],
  "input capture": ["T1056"],
  "credential access": ["T1003", "T1555"],
  "token manipulation": ["T1134"],
  "session hijacking": ["T1563.001"],
  "cookie theft": ["T1539"],
  "browser exploitation": ["T1189"],
  "drive-by compromise": ["T1189"],
  "watering hole": ["T1189"],
  "supply chain attack": ["T1195"],
  "living off the land": ["T1218"],
  "dual use tools": ["T1218"],
  "scripting": ["T1059"],
  "PowerShell": ["T1059.001"],
  "Python": ["T1059.006"],
  "Bash": ["T1059.004"],
  "JavaScript": ["T1059.007"],
  "VBScript": ["T1059.005"],
  "macro": ["T1059.005"],
  "executable": ["T1204.002"],
  "DLL": ["T1204.002"],
  "service": ["T1543.003"],
  "scheduled task": ["T1053.005"],
  "startup": ["T1547.001"],
  "registry run": ["T1547.001"],
  "WMI": ["T1047"],
  "RDP": ["T1021.001"],
  "SSH": ["T1021.004"],
  "VPN": ["T1071.001"],
  "Tor": ["T1071.001"],
  "proxy": ["T1090"],
  "tunnel": ["T1572"],
  "DNS tunneling": ["T1071.004"],
  "ICMP tunneling": ["T1071.003"],
  "HTTP tunneling": ["T1071.001"],
  "HTTPS tunneling": ["T1071.001"],
  "SMB": ["T1021.002"],
  "FTP": ["T1071.002"],
  "SMTP": ["T1071.003"],
  "IMAP": ["T1071.003"],
  "POP3": ["T1071.003"],
  "RPC": ["T1021.002"],
  "LDAP": ["T1087.002"],
  "Kerberos": ["T1558.003"],
  "NTLM": ["T1558.002"],
  "Active Directory": ["T1087.002"],
  "Group Policy": ["T1484.001"]
}